# deploy/gpn-web.service
[Unit]
Description=GPN web (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gpn
WorkingDirectory=/opt/gpn/web
# NEXT_PUBLIC_SIGNALING_URL is baked in at build time, not needed here.
ExecStart=/opt/gpn/web/node_modules/.bin/next start -p 3000
Restart=on-failure
RestartSec=3

Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
